# system/skeleton/etc/network/if-up.d/check_ip_assigned
#!/bin/sh
IP_ASSIGNMENT_FILE=/var/run/ip_assign_method

assigned_by() {
  method=$1
  echo "$method" > $IP_ASSIGNMENT_FILE
}

if [ "$LOGICAL" = "wlan0" ]; then
  assigned_by "$METHOD"
fi
